# infra/mission-control.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.quartermint.mission-control</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/services/mission-control/node_modules/.bin/tsx</string>
        <string>packages/api/src/index.ts</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/opt/services/mission-control</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>NODE_ENV</key>
        <string>production</string>
        <key>PORT</key>
        <string>3000</string>
        <key>MC_CONFIG_PATH</key>
        <string>/opt/services/mission-control/mc.config.json</string>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/opt/services/mission-control/logs/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/opt/services/mission-control/logs/stderr.log</string>

    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>4096</integer>
    </dict>
</dict>
</plist>
